# go-firewall/go-firewall.service
[Unit]
Description=Go Firewall Process Monitor
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/go-firewall
ExecStart=/opt/go-firewall/go-firewall
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
